# terraform/postgres-backup/postgres-backup.timer
[Unit]
Description=Runs PostgreSQL Backup Service Hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
